# misc/docker/supervisord.ini
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB   ; maximum size of logfile before rotation
logfile_backups=2       ; number of backed up logfiles
loglevel=info           ; info, debug, warn, trace

[program:uwsgi]
command=/usr/bin/uwsgi --uid blunders --gid=blunders /home/blunders/uwsgi.ini
stdout_logfile=/var/log/uwsgi.log
stderr_logfile=/var/log/uwsgi.log

[program:redis]
command=/usr/bin/redis-server /etc/redis.conf

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf -g "daemon off;"

[unix_http_server]
username=root
password=chessdb
