# misc/docker/supervisord.ini
[supervisord]
nodaemon=true

[program:uwsgi]
command=/usr/bin/uwsgi --uid blunders --gid=blunders /home/blunders/uwsgi.ini

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf -g "daemon off;"
